# src/etc/supervisord.d/mysqld-bootstrap.conf
[program:mysqld-bootstrap]
priority = 6
command = /usr/sbin/mysqld-bootstrap
autostart = %(ENV_MYSQL_AUTOSTART_MYSQLD_BOOTSTRAP)s
startsecs = 0
startretries = 0
autorestart = false
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
