# etc/services-config/supervisor/supervisord.d/mysqld-wrapper.conf
[program:mysqld-wrapper]
priority = 100
command = /usr/sbin/mysqld-wrapper
startsecs = 0
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/mysqld.log
stdout_events_enabled = true