# etc/services-config/supervisor/supervisord.d/mysqld-bootstrap.conf
[program:mysqld-bootstrap]
priority = 6
command = bash -c 'touch /tmp/mysqld-bootstrap.lock; /usr/sbin/mysqld-bootstrap && rm -f /tmp/mysqld-bootstrap.lock'
startsecs = 0
startretries = 0
autorestart = false
redirect_stderr = true
stdout_logfile = /var/log/mysqld.log
stdout_events_enabled = true